# claude-telegram/claude-telegram-restart.service
[Unit]
Description=Restart claude-telegram (picks up claude update + fresh session)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart claude-telegram.service
